# docker/app/supervisor.conf
[supervisord]
nodaemon=true
[program:app-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/artisan queue:work --sleep=3 --tries=3
autostart=trueautorestart=truenumprocs=2
redirect_stderr=trueuser=root
stdout_logfile=/var/www/app/storage/logs/worker.log
